# mashpad.service
[Unit]
Description=mashpad — baby keyboard smasher
After=systemd-user-sessions.service sound.target seatd.service
Wants=seatd.service
Conflicts=getty@tty1.service

[Service]
User=@USER@
# _seatd: reach the seatd socket so cage can acquire the seat (DRM + VT).
SupplementaryGroups=video input audio render tty _seatd
WorkingDirectory=@REPO@
Environment=PYTHONPATH=@REPO@
# cage is a single-app Wayland kiosk. mashpad renders inside it (via XWayland,
# SDL driver "x11") because SDL's bare kmsdrm backend does not present on the
# Pi's split-GPU (v3d renders, vc4 scans out). SDL_VIDEODRIVER is left UNSET on
# purpose — under cage, SDL's default selection is the path that renders.
Environment=XDG_RUNTIME_DIR=/run/mashpad
RuntimeDirectory=mashpad
RuntimeDirectoryMode=0700
ExecStart=/usr/bin/cage -- /usr/bin/python3 -m mashpad
Restart=always
RestartSec=2
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=journal

[Install]
WantedBy=multi-user.target
